# files/lib/systemd/reflector.timer
[Unit]
Description=Weekly pacman mirrorlist refresh

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=12h

[Install]
WantedBy=timers.target
